# deployment/supervisor/microblog-tasks.conf
[program:foody-tasks]
command=/home/ubuntu/foody/venv/bin/rq worker foody-tasks
numprocs=1
directory=/home/ubuntu/foody
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
